# root/etc/services.d/mfi/run
#!/usr/bin/with-contenv bash

cd /config

exec /usr/bin/java -jar /usr/lib/mfi/lib/ace.jar start
