# root/etc/services.d/mfi/run
#!/usr/bin/with-contenv bash
cd /config
exec s6-setuidgid abc java -Xmx1024M -jar /mFi/lib/ace.jar start
