# root/etc/services.d/mfi/run
#!/usr/bin/with-contenv bash

#cd /config
cd /usr/lib/mfi/
exec s6-setuidgid abc java -Xmx1024M -jar /usr/lib/mfi/lib/ace.jar start
